# systemd/reboot.service
[Unit]
Description=Scheduled system reboot
Wants=reboot.timer

[Service]
Type=oneshot
ExecStart=/sbin/reboot
